# infra/docker/sch/crontab/root
# min hour day month weekday command
# * * * * * /bin/date --rfc-2822 >> /var/log/cron.log
* * * * * node /home/node/app/scheduled.js >> /var/log/cron.log 2>&1
